# assets/ddns.service
[Unit]
Description=DDNS Route53 Service

[Service]
Type=oneshot
ExecStart=/usr/bin/env ddns-route53 update
